# script-plus-conf/conf.d/phpd.conf
[program:phpd]
command= /usr/sbin/php-fpm7.4 -F
autorestart=true
