# script-plus-conf/conf.d/phpd.conf
[program:phpd]
command= php-fpm7 --nodaemonize
autorestart=true
