# script-plus-conf/conf.d/phpd.conf
[program:phpd]
command= php-fpm -F
autorestart=true
